# data/pull-metrics.timer
[Unit]
Description=Pull metrics repo (timer)
After=timers.target

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target