# data/start-all-timers.timer
[Unit]
Description=Start all metric collectors (timer)
After=timers.target

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target